# conf.d/django.conf
[program:django]
command=gunicorn optify.wsgi:application
directory=/optify
autostart=true
autorestart=true
