# conf.d/tornado.conf
[program:tornado]
command= python3.8 index.py
directory = /optify/chat/tornado_server
autostart=true
autorestart=true
